# packages/hotpi/system.d/adaway.timer
[Unit]
# The description should be used to explain what this servicefile is for
Description=Runs Adaway every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
