# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon
ConditionKernelCommandLine=|ssh
ConditionPathExists=|/storage/.cache/services/sshd.conf
After=network.target

[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
ExecStart=/usr/sbin/sshd -D -f /storage/.cache/services/sshd.conf
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
