# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon
After=network.target sshd-defaults.service
Requires=sshd-defaults.service

ConditionPathExists=/storage/.cache/services/sshd.conf


[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
ExecStart=/usr/sbin/sshd -D -f /storage/.cache/services/sshd.conf
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
